# slurm-login/update-motd.d/00-header
printf "Welcome to a CoreWeave SLURM HPC Cluster"